# crontab
* * * * * /home/simon/txt >/tmp/txt 2>&1
